# systemd/fang@.service
[Unit]
Description=fang agent (%i)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/fang
ExecStart=node --experimental-strip-types src/cli.ts wrap %i --port 3001
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
